# deploy/systemd/artctl.service
[Unit]
Description=ARTCTL API service
After=network.target

[Service]
Type=simple
User=cortex
Group=cortex
WorkingDirectory=/srv/artctl
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
